# crius.service
[Unit]
Description=crius CRI container runtime
Wants=network-online.target
After=network-online.target dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/crius --config /etc/crius/crius.conf
Restart=always
RestartSec=5
Delegate=yes
KillMode=process
StateDirectory=crius
LogsDirectory=crius

# Keep process limits out of the daemon accounting path; use cgroups for
# container-local accounting.
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
OOMScoreAdjust=-999

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
